# files/systemd/system/getty@.service.d/autologin.conf
# https://wiki.archlinux.org/title/Getty#Automatic_login_to_virtual_console

[Service]
ExecStart=
ExecStart=-/sbin/agetty --noreset --autologin bleroux - ${TERM}
