# files/systemd/system/battery-charge-limit.service
[Unit]
Description=Set Battery Charge Limit
After=multi-user.targe
ConditionPathExists=/sys/class/power_supply/BAT0/charge_control_end_threshold

[Service]
Type=oneshot
ExecStart=/bin/sh -c 'echo 80 > /sys/class/power_supply/BAT0/charge_control_end_threshold'

[Install]
WantedBy=multi-user.target
