# files/systemd/user/niri.service.d/autostart.conf
# Copied with bluebuild

# Added
[Install]
WantedBy=default.target
